# bookscout.service
# BookScout systemd service file (API)
#
# NOTE: this unit runs only the API. Scheduled scans, imports, and every
# queued job are executed by the arq worker — install bookscout-worker.service
# alongside this one or no background job will ever run.
#
# Installation:
# 1. Copy this file to /etc/systemd/system/bookscout.service
#    (and bookscout-worker.service to /etc/systemd/system/)
# 2. Update paths and user below
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable bookscout bookscout-worker
# 5. Run: sudo systemctl start bookscout bookscout-worker

[Unit]
Description=BookScout - Multi-source book discovery
After=network.target

[Service]
Type=simple
User=your-username
Group=your-username
WorkingDirectory=/path/to/bookscout
Environment="PATH=/path/to/venv/bin"
# Apply pending DB migrations before starting the API
ExecStartPre=/path/to/venv/bin/alembic upgrade head
ExecStart=/path/to/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8765
Restart=always
RestartSec=10

# Optional: Environment variables
# Environment="AUDIOBOOKSHELF_URL=http://localhost:13378"
# Environment="AUDIOBOOKSHELF_TOKEN=your_token"
# Environment="PROWLARR_URL=http://localhost:9696"
# Environment="PROWLARR_API_KEY=your_key"

[Install]
WantedBy=multi-user.target
